# srcs/mysql/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:mysql]
command=mysql
autostart=true
autorestart=true