# srcs/grafana/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:grafana]
directory=/usr/share/grafana/conf
command= grafana-server --homepath=/usr/share/grafana
autostart=true
autorestart=true
sdout_logfile=/var/log/supervisor/%(program_name)sdout.log
sderr_logfile=/var/log/supervisor/%(program_name)sderr.log

[program:telegraf]
command= telegraf --config /etc/telegraf/telegraf.conf
autostart=true
autorestart=true